# docker/supervisord/supervisord.octane.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777
# chown=supervisor:supervisor

[supervisord]
nodaemon=true
#user=supervisor
stdout_logfile=/var/log/supervisor/supervisord.log

[program:nginx]
command=nginx -c /etc/nginx/nginx.conf -g 'daemon off;'
autostart=true
user=root
#user=supervisor

[program:laravel-octane]
command=docker-php-entrypoint php -d variables_order=EGPCS /var/www/app/artisan octane:start --server=swoole --host=0.0.0.0 --port=8000
autostart=true
autorestart=true
user=root
#user=supervisor

[program:horizon]
process_name=%(program_name)s
command=docker-php-entrypoint php /var/www/app/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
user=root
#user=supervisor

[program:horizon-workers-high]
process_name=%(program_name)s_%(process_num)s
command=docker-php-entrypoint php /var/www/app/artisan queue:work --queue=high --sleep=3 --tries=3
autostart=true
autorestart=true
redirect_stderr=true
user=root
numprocs=10

[program:horizon-workers-default]
process_name=%(program_name)s_%(process_num)s
command=docker-php-entrypoint php /var/www/app/artisan queue:work --queue=default --sleep=3 --tries=3
autostart=true
autorestart=true
redirect_stderr=true
user=root
numprocs=10

[program:horizon-workers-low]
process_name=%(program_name)s_%(process_num)s
command=docker-php-entrypoint php /var/www/app/artisan queue:work --queue=low --sleep=3 --tries=3
autostart=true
autorestart=true
redirect_stderr=true
user=root
numprocs=10